# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
silent=false

; ── Tidarr (Tidal downloader) — interne poort 8484 ──────────────────────────
[program:tidarr]
command=sh /tidarr/docker/entrypoint.sh
directory=/tidarr
autostart=true
autorestart=true
priority=10
startsecs=10
stopasgroup=true
stopwaittimeout=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
environment=HOME="/shared",ENVIRONMENT="production",PORT="8484"

; ── MediaSage (LLM muziekaanbevelingen) — interne poort 5765 ────────────────
[program:mediasage]
command=sh -c "mkdir -p /mediasage/data && chmod 777 /mediasage/data && python -m uvicorn backend.main:app --host 0.0.0.0 --port 5765 --workers 1"
directory=/mediasage
autostart=true
autorestart=true
priority=15
startsecs=10
stopasgroup=true
stopwaittimeout=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
environment=PYTHONUNBUFFERED="1"

; ── Muziekdashboard (lastfm-app) — poort 80 ────────────────────────────────
[program:lastfm]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=5
stopasgroup=true
stopwaittimeout=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
environment=NODE_ENV="production",LOG_LEVEL="info"

; ── AudioMuse Flask webserver — interne poort (bijv. 5000) ──────────────────
[program:audiomuse-flask]
command=/app/venv/bin/python /app/audiomuse/app.py
directory=/app/audiomuse
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=5
stopasgroup=true
stopwaittimeout=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
environment=PYTHONUNBUFFERED="1",MEDIASERVER_TYPE="plex",REDIS_URL="redis://redis:6379/0",POSTGRES_HOST="postgres",DATABASE_URL="%(ENV_DATABASE_URL)s",PLEX_URL="%(ENV_PLEX_URL)s",PLEX_TOKEN="%(ENV_PLEX_TOKEN)s",SECRET_KEY="%(ENV_SECRET_KEY)s",FLASK_ENV="production",AUTH_ENABLED="true",AUDIOMUSE_USER="%(ENV_AUDIOMUSE_USER)s",AUDIOMUSE_PASSWORD="%(ENV_AUDIOMUSE_PASSWORD)s",ENABLE_PROXY_FIX="true"

; ── AudioMuse RQ Worker — audio-analyse taken ────────────────────────────────
[program:audiomuse-worker]
command=/app/venv/bin/python /app/audiomuse/rq_worker.py
directory=/app/audiomuse
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=5
stopasgroup=true
stopwaittimeout=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
environment=PYTHONUNBUFFERED="1",AUDIOMUSE_ROLE="worker",MEDIASERVER_TYPE="plex",REDIS_URL="redis://redis:6379/0",POSTGRES_HOST="postgres",DATABASE_URL="%(ENV_DATABASE_URL)s",PLEX_URL="%(ENV_PLEX_URL)s",PLEX_TOKEN="%(ENV_PLEX_TOKEN)s",SECRET_KEY="%(ENV_SECRET_KEY)s",FLASK_ENV="production",AUTH_ENABLED="true",AUDIOMUSE_USER="%(ENV_AUDIOMUSE_USER)s",AUDIOMUSE_PASSWORD="%(ENV_AUDIOMUSE_PASSWORD)s"

; ── Event listeners voor process state changes ───────────────────────────────
[eventlistener:processes]
buffer_size=100
events=PROCESS_STATE
command=/bin/bash -c "while read line; do echo \"[$(date '+%%Y-%%m-%%d %%H:%%M:%%S')] [supervisord] $${line}\"; done"
