# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Tidarr (Tidal downloader) — interne poort 8484 ──────────────────────────
[program:tidarr]
command=sh /tidarr/docker/entrypoint.sh
directory=/tidarr
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/shared",ENVIRONMENT="production"

; ── Muziekdashboard (lastfm-app) — poort 80 ────────────────────────────────
[program:lastfm]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
